# src/roles/apache_airflow/templates/airflow-scheduler.service.j2
[Unit]
Description=Apache Airflow Scheduler
After=network.target
Wants=network.target

[Service]
EnvironmentFile=/etc/default/airflow
User={{ airflow_user }}
Group={{ airflow_group }}
Type=simple
ExecStart={{ airflow_venv_path | default('') }}/bin/airflow scheduler
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
