# src/roles/infrastructure/apt_mirror/templates/apt-mirror-alerts.service.j2
[Unit]
Description=APT Mirror Alert Manager Service
After=network.target

[Service]
Type=oneshot
User=apt-mirror
Group=apt-mirror
WorkingDirectory=/opt/apt-mirror-alerts
ExecStart=/usr/bin/python3 /opt/apt-mirror-alerts/alert_manager.py
Environment=PYTHONPATH=/opt/apt-mirror-alerts

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/apt-mirror-alerts /var/log /var/lib/apt-mirror-alerts
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
