# src/roles/data_analytics/spark_role/templates/spark-master.service.j2
[Unit]
Description=Apache Spark Master
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User={{ spark_role_user }}
Group={{ spark_role_group }}
WorkingDirectory={{ spark_role_symlink_dir }}/sbin
ExecStart={{ spark_role_symlink_dir }}/sbin/start-master.sh
ExecStop={{ spark_role_symlink_dir }}/sbin/stop-master.sh
{% if spark_role_health_checks_enabled %}
ExecStartPost=/bin/sleep 30
ExecStartPost={{ spark_role_install_dir }}/scripts/health_check.sh
{% endif %}
SuccessExitStatus=143
Restart=on-failure
RestartSec=30
{% if spark_role_health_checks_enabled %}
# Health check configuration
TimeoutStartSec=300
TimeoutStopSec=120
{% endif %}

[Install]
WantedBy=multi-user.target
