# src/roles/spark_role/templates/spark-master.service.j2
[Unit]
Description=Apache Spark Master
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User={{ spark_user }}
Group={{ spark_group }}
WorkingDirectory={{ spark_symlink_dir }}/sbin
ExecStart={{ spark_symlink_dir }}/sbin/start-master.sh
ExecStop={{ spark_symlink_dir }}/sbin/stop-master.sh
SuccessExitStatus=143
Restart=on-failure

[Install]
WantedBy=multi-user.target
