# src/roles/data_analytics/apache_superset/templates/superset-beat.service.j2
[Unit]
Description=Superset Celery Beat Scheduler
After=network.target
Requires={{ superset_service_name }}.service

[Service]
Type=simple
User={{ superset_user }}
Group={{ superset_group }}
WorkingDirectory={{ superset_install_dir }}
Environment=SUPERSET_CONFIG_PATH={{ superset_config_path }}
ExecStart={{ superset_venv_dir }}/bin/celery beat --app=superset.tasks.celery_app:app --loglevel={{ superset_celery_worker_log_level }}
Restart=always
RestartSec=10

# Logging
StandardOutput=file:/var/log/celery/{{ superset_service_name }}-beat.log
StandardError=file:/var/log/celery/{{ superset_service_name }}-beat-error.log

# Process management
PIDFile=/var/run/celery/{{ superset_service_name }}-beat.pid
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ superset_install_dir }} /var/log/celery /var/run/celery

[Install]
WantedBy=multi-user.target
