# src/roles/data_analytics/apache_superset/templates/superset-worker.service.j2
[Unit]
Description=Apache Superset Celery Worker
After=network.target
Requires={{ apache_superset_service_name }}.service

[Service]
Type=simple
User={{ apache_superset_user }}
Group={{ apache_superset_group }}
WorkingDirectory={{ apache_superset_install_dir }}
Environment=SUPERSET_CONFIG_PATH={{ apache_superset_config_path }}
ExecStart={{ apache_superset_venv_dir }}/bin/celery --app=superset.tasks.celery_app:app worker --pool=prefork -O fair -c {{ apache_superset_celery_worker_concurrency }} --loglevel={{ apache_superset_celery_worker_log_level }}
Restart=always
RestartSec=10

# Logging
StandardOutput=file:/var/log/celery/{{ apache_superset_service_name }}-worker.log
StandardError=file:/var/log/celery/{{ apache_superset_service_name }}-worker-error.log

# Process management
PIDFile=/var/run/celery/{{ apache_superset_service_name }}-worker.pid
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ apache_superset_install_dir }} /var/log/celery /var/run/celery

[Install]
WantedBy=multi-user.target
