# src/roles/data_analytics/apache_superset/templates/superset.service.j2
[Unit]
Description=Apache Superset
After=network.target

[Service]
Type=simple
User={{ apache_superset_user }}
Group={{ apache_superset_group }}
WorkingDirectory={{ apache_superset_install_dir }}
Environment=SUPERSET_CONFIG_PATH={{ apache_superset_config_path }}
{% if apache_superset_ssl_enabled %}
ExecStart={{ apache_superset_venv_dir }}/bin/gunicorn \
    -w {{ apache_superset_gunicorn_workers }} \
    -k {{ apache_superset_gunicorn_worker_class }} \
    --timeout {{ apache_superset_gunicorn_timeout }} \
    --keepalive {{ apache_superset_gunicorn_keepalive }} \
    --max-requests {{ apache_superset_gunicorn_max_requests }} \
    --max-requests-jitter {{ apache_superset_gunicorn_max_requests_jitter }} \
    {% if apache_superset_gunicorn_preload_app %}--preload{% endif %} \
    --worker-connections {{ apache_superset_gunicorn_worker_connections }} \
    -b {{ apache_superset_gunicorn_bind }} \
    --certfile {{ apache_superset_ssl_cert_path }} \
    --keyfile {{ apache_superset_ssl_key_path }} \
    {% if apache_superset_ssl_ca_path %}--ca-certs {{ apache_superset_ssl_ca_path }}{% endif %} \
    --ssl-version TLSv1_2 \
    --ciphers {{ apache_superset_ssl_ciphers }} \
    --limit-request-line 0 \
    --limit-request-field_size 0 \
    {% if apache_superset_enable_access_log %}--access-logfile {{ apache_superset_access_log_file }}{% endif %} \
    --log-file {{ apache_superset_log_file }} \
    --log-level {{ apache_superset_log_level.lower() }} \
    superset.app:create_app()
{% else %}
ExecStart={{ apache_superset_venv_dir }}/bin/gunicorn \
    -w {{ apache_superset_gunicorn_workers }} \
    -k {{ apache_superset_gunicorn_worker_class }} \
    --timeout {{ apache_superset_gunicorn_timeout }} \
    --keepalive {{ apache_superset_gunicorn_keepalive }} \
    --max-requests {{ apache_superset_gunicorn_max_requests }} \
    --max-requests-jitter {{ apache_superset_gunicorn_max_requests_jitter }} \
    {% if apache_superset_gunicorn_preload_app %}--preload{% endif %} \
    --worker-connections {{ apache_superset_gunicorn_worker_connections }} \
    -b {{ apache_superset_gunicorn_bind }} \
    --limit-request-line 0 \
    --limit-request-field_size 0 \
    {% if apache_superset_enable_access_log %}--access-logfile {{ apache_superset_access_log_file }}{% endif %} \
    --log-file {{ apache_superset_log_file }} \
    --log-level {{ apache_superset_log_level.lower() }} \
    superset.app:create_app()
{% endif %}
Restart=always
RestartSec=10

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ apache_superset_install_dir }} {{ apache_superset_log_file | dirname }} {% if apache_superset_enable_access_log %}{{ apache_superset_access_log_file | dirname }}{% endif %}

[Install]
WantedBy=multi-user.target
