# src/databases/linux/roles/postgres_backup/templates/postgres_backup.service.j2
[Unit]
Description=PostgreSQL database backup service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_postgres.sh
