# src/databases/linux/roles/mariadb_backups/files/mariadb-backup.service
[Unit]
Description=MariaDB Backup

[Service]
ExecStart=/usr/local/bin/mariadb-backup.sh

[Install]
WantedBy=multi-user.target
